# roles/minecraft/templates/overviewer.cron
0 * * * *   minecraft_overviewer    cd /usr/local/Minecraft-Overviewer && ./overviewer.py -c /etc/minecraft-overviewer.yml
